# api/supervisord/stuhlgang.conf
[unix_http_server]
file=/tmp/stuhlgang-supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/stuhlgang-supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/stuhlgang-supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/stuhlgang-supervisor.sock ; use a unix:// URL  for a unix socket

[program:webapp]
command= /home/matt/.virtualenvs/stuhlgang/bin/stuhlgang-run-webapp dev.yaml
environment=HOME="/home/matt",USER="matt",VIRTUAL_ENV="$/home/matt/.virtualenvs/stuhlgang"

[eventlistener:crashmail]
command=/home/matt/.virtualenvs/stuhlgang/bin/crashmail -a -m info@example.com
events=PROCESS_STATE

[eventlistener:memmon]
command=/home/matt/.virtualenvs/stuhlgang/bin/memmon -a 1500MB -m info@example.com
events=TICK_60
